# devops-environment-meta/debian/postinst
#!/bin/sh

invoke-rc.d libvirt-bin restart


