# opentracker/debian/postrm
#!/bin/bash

if [[ -x /etc/init.d/opentracker ]] ; then
    invoke-rc.d opentracker stop
fi

